# cron/crontab
*/2 * * * * php bin/console parse:news >> /var/log/cron/cron.log 2>&1